# web-message-queue.service
[Unit]
Description=web message queue
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/web
ExecStart=/usr/local/bin/celery -A tasks worker --loglevel=info
Restart=always

[Install]
WantedBy=multi-user.target